# docker/supervisor/conf.d/25-vite.conf
[program:vite]
command=sh -c "export DEV_PROXY_TARGET=http://127.0.0.1:8000 && npm install && exec npm run dev -- --host 0.0.0.0 --port 5173"
directory=/app/frontend
autostart=true
autorestart=unexpected
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
